# emulator/system/vcom.target
[Unit]
Description="Virtual Ports"
Requires=vcom@1.service vcom@2.service

[Install]
WantedBy=multi-user.target
